# files/etc/uci-defaults/10_set_russian_lang
#!/bin/sh
uci set luci.main.lang='ru'
uci commit luci
rm /etc/uci-defaults/10_set_russian_lang
exit 0
