# files/etc/uci-defaults/99_sysupgrade_fix
#!/bin/sh
cat > /etc/sysupgrade.conf << EOF
CHECK_DEVICE=0
FORCE=1
KEEP_CONFIG=1
NO_CHECK_MD5=1
NO_CHECK_SIGNATURE=1
ALLOW_BACKUP=1
EOF

rm -f /etc/uci-defaults/99_sysupgrade_fix
exit 0
